# supervisord.conf
[supervisord]
nodaemon=true

[program:telldus]
command=/etc/init.d/telldusd start

[program:nodered]
command=/usr/local/bin/node-red-pi -v --max-old-space-size=128
